# conf/etc/systemd/system/systemd-hibernate.service.d/maintenance.conf
# See systemd-suspend.service.d/maintenance.conf for rationale.
[Service]
ExecStartPost=/etc/systemd/system-sleep/50-user-maintenance post hibernate
