# conf/home/systemd/user/maintenance.service
[Unit]
Description=User maintenance (flatpak updates, etc.)
# Fired at login via default.target, and on resume via the system-sleep hook
# installed under /etc/systemd/system-sleep/. The script throttles itself,
# so duplicate triggers are harmless.

[Service]
Type=oneshot
ExecStart=%h/.local/bin/maintenance.zsh
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=30min
# Send stdout/stderr to the journal under a stable identifier so the run can be
# followed live or after the fact:
#   journalctl --user -t maintenance        # this unit's output, all runs
#   journalctl --user -u maintenance -f     # follow the unit live
StandardOutput=journal
StandardError=journal
SyslogIdentifier=maintenance

[Install]
WantedBy=default.target
