# conf/home/systemd/user/maintenance.service
[Unit]
Description=User maintenance (flatpak updates, etc.)
# Fired at login via default.target, and on resume via the system-sleep hook
# installed under /etc/systemd/system-sleep/. The script throttles itself,
# so duplicate triggers are harmless.

[Service]
Type=oneshot
ExecStart=%h/.local/bin/maintenance.zsh
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=30min

[Install]
WantedBy=default.target
