# conf/etc/systemd/system/systemd-suspend.service.d/maintenance.conf
# systemd-sleep(8) only auto-runs hooks from /usr/lib/systemd/system-sleep/.
# Rather than touch /usr/lib, hook into the sleep service's ExecStartPost,
# which fires after the kernel resumes (the main process is what performs
# the suspend and only returns on wake).
[Service]
ExecStartPost=/etc/systemd/system-sleep/50-user-maintenance post suspend
